# ansible/roles/base_24/files/motd/update-motd.d/90-pemmican
#!/bin/sh

if [ -x /usr/bin/pemmican-cli ]; then
    exec /usr/bin/pemmican-cli
fi
